# etc/sd-proxy@.socket
[Socket]
ListenStream=8080
ReusePort=true
NoDelay=true
IPTOS=low-delay
